# deploy/systemd/election-api.service
[Unit]
Description=Election Manager API (FastAPI)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/election-manager/backend

# Use the venv Python if present (recommended)
# ExecStart=/opt/election-manager/.venv/bin/gunicorn -k uvicorn.workers.UvicornWorker -w 4 -b 127.0.0.1:8000 main:app

# Fallback (system python) - ensure gunicorn is installed in that env
ExecStart=/usr/bin/python3 -m gunicorn -k uvicorn.workers.UvicornWorker -w 4 -b 127.0.0.1:8000 main:app

Restart=always
RestartSec=3
Environment=PYTHONUNBUFFERED=1

# DB/SQLite tuning (if you are still on SQLite; prefer PostgreSQL for high load)
Environment=SQLITE_BUSY_TIMEOUT_SEC=30

[Install]
WantedBy=multi-user.target
